# archlinux-anvl/mwifiex-p2p/70-mwifiex-p2p.rules
KERNEL=="mlan0", ACTION=="add", RUN+="/usr/bin/iw phy phy0 interface add p2p0 type __p2pcl"

